# config/podman-prune.timer
[Unit]
Description=Weekly Podman Cleanup

[Timer]
OnCalendar=Sun 03:00
Persistent=true

[Install]
WantedBy=timers.target
